# po8e/PO8e-for-3.2.0-3-amd64/Plx9054.rules
KERNEL=="Plx9054", SUBSYSTEM=="module", RUN+="Plx9054.sh"
